# gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn offsetter:app -b localhost:5000 -w 4
directory=/deploy/app
